# crontab
#每周一0点1分执行
1 0 * * 1 /opt/rsq_s5rt/rsq_s5rt_start.sh